# aqi-train-online-ar.service
[Unit]
Description=run AQPy adaptive AR online training step
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/true
WorkingDirectory=/home/pi/AQPy
EnvironmentFile=-/home/pi/AQPy/.env
Environment=PYTHONDONTWRITEBYTECODE=1
User=pi
Group=pi
NoNewPrivileges=true
PrivateTmp=true
PrivateMounts=true
ProtectSystem=full
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
SystemCallArchitectures=native
UMask=0077
