# aqi-forecast.service
[Unit]
Description=run AQPy forecast inference
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/pi/AQPy
EnvironmentFile=-/home/pi/AQPy/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/pi/AQPy/.venv/bin/python /home/pi/AQPy/run_forecast_batch.py --spec-file /home/pi/AQPy/configs/model_specs.json
User=pi
Group=pi
NoNewPrivileges=true
PrivateTmp=true
PrivateMounts=true
ProtectSystem=full
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
SystemCallArchitectures=native
UMask=0077
