# aqi-train-online.service
[Unit]
Description=run AQPy online NN training step
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/pi/AQPy
EnvironmentFile=-/home/pi/AQPy/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/pi/AQPy/scripts/run_edge_jobs_now.sh --train-only
User=pi
Group=pi
NoNewPrivileges=true
PrivateTmp=true
PrivateMounts=true
ProtectSystem=full
ProtectHome=false
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
SystemCallArchitectures=native
ReadWritePaths=/home/pi/AQPy/models
UMask=0077
