# aqi.service
[Unit]
Description=run AQI sensor script

[Service]
Type=simple
WorkingDirectory=/home/pi/AQPy
EnvironmentFile=-/home/pi/AQPy/.env
Environment=PYTHONDONTWRITEBYTECODE=1
# use -u to make print statements unbuffered in logs
ExecStart=/home/pi/AQPy/.venv/bin/python -u /home/pi/AQPy/read_sensors.py
User=pi
Group=pi
Restart=on-failure
RestartSec=10s
NoNewPrivileges=true
PrivateTmp=true
PrivateMounts=true
ProtectSystem=full
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
